# vps-logger-bot.service.template
[Unit]
Description=Telegram auth.log bot
After=network.target

[Service]
# Path to the file to be used for logging
Environment=LOGFILE_PATH=_LOGFILE_PATH
# Bot token provided by Bot Father
Environment=BOT_TOKEN=_BOT_TOKEN
Type=simple
User=root
WorkingDirectory=/home/someuser
ExecStart=/snap/bin/node /path/to/start.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target